# sdk/examples/listener_os/config/activation/linux/listener-os-dispatch.path.tmpl
[Unit]
Description=Watch listener_os trigger file

[Path]
PathChanged={{TRIGGER_FILE}}
Unit={{SERVICE_NAME}}.service

[Install]
WantedBy=default.target
