# sdk/examples/listener_os/config/activation/linux/listener-os-dispatch.service.tmpl
[Unit]
Description=Listener OS dispatch once

[Service]
Type=oneshot
WorkingDirectory={{WORKDIR}}
ExecStart={{PYTHON}} -m listener_os.main dispatch-once --db-path {{DB_PATH}}
